# systemd/poly-collector.service
[Unit]
Description= Polymarket Data Collector
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/poly-edgefinder
Environment="PATH=/opt/poly-edgefinder/venv/bin"
EnvironmentFile=/opt/poly-edgefinder/.env
ExecStart=/opt/poly-edgefinder/venv/bin/python -m src.data_collectors.simmer_collector
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target